# systemd/brew-bootstrap.service
[Unit]
Description=Install default Homebrew formulas
Wants=network-online.target
After=network-online.target brew-setup.service
# Run once only; drop a stamp when complete
ConditionPathExists=!/var/lib/ublue/brew-bootstrap.stamp

[Service]
Type=oneshot
# Run as the current user to fix linuxbrew user does not exist
User=%u
Group=%u
# Brew needs a real HOME and a readable Brewfile
Environment=HOME=/var/home/linuxbrew
Environment=NONINTERACTIVE=1
Environment=HOMEBREW_NO_ENV_HINTS=1
Environment=BASH_ENV=/etc/profile
WorkingDirectory=/var/home/linuxbrew
ExecStart=/usr/bin/bash -lc '\
  eval "$(/var/home/linuxbrew/.linuxbrew/bin/brew shellenv)"; \
  /var/home/linuxbrew/.linuxbrew/bin/brew update --auto-update || true; \
  /var/home/linuxbrew/.linuxbrew/bin/brew bundle --no-lock --file=/etc/ublue/Brewfile \
'
ExecStartPost=/usr/bin/install -D -m 0644 /dev/null /var/lib/ublue/brew-bootstrap.stamp

[Install]
WantedBy=multi-user.target

