# files/brew-bootstrap.service
[Unit]
Description=Dummy placeholder for brew bootstrap

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
WantedBy=default.target
