# files/brew-once.service
[Unit]
Description=Run brew bundle once on first login
ConditionPathExists=%h/.Brewfile

[Service]
Type=oneshot
Environment=NONINTERACTIVE=1
ExecStartPre=/usr/bin/echo ">>> Installing Homebrew packages from ~/.Brewfile, this may take a few minutes..." > /dev/tty
ExecStart=/usr/bin/bash -lc '\
  set -euo pipefail; \
  BREW="/home/linuxbrew/.linuxbrew/bin/brew"; \
  [ -x "$BREW" ] || BREW="/var/home/linuxbrew/.linuxbrew/bin/brew"; \
  if [ ! -x "$BREW" ]; then \
    echo "Error: brew not found" >&2; exit 1; \
  fi; \
  echo ">>> Installing Homebrew packages from ~/.Brewfile..."; \
  "$BREW" bundle --file=%h/.Brewfile || true; \
  echo ">>> Brewfile installation finished. Cleaning up..."; \
  rm -f %h/.Brewfile \
'

[Install]
WantedBy=default.target
